# contrib/runit/log/run
#!/bin/sh
mkdir -p /var/log/lokinet
exec chpst -u root:root svlogd -tt /var/log/lokinet/
